# .config/autostart/dunst.desktop
[Desktop Entry]
Type=Application
Name=Dunst
Exec=dunst
X-GNOME-Autostart-enabled=true
